# PetaLinux/FinalProject/project-spec/meta-system-controller/recipes-utils/bootbin-fw-update/files/bootbin-fw-update.service
[Unit]
Description=BOOT.BIN Update Status
After=basic.target

[Service]
ExecStart=/usr/bin/bootbin-fw-update.sh
ExecStartPre=/bin/echo "BOOT.BIN Update Status"
Type=oneshot

[Install]
WantedBy=multi-user.target
